# monitor.service
[Unit]
Description=Air monitor
[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/python3 /home/james-monitor/demo.py
[Install]
WantedBy=multi-user.target
